# deploy/health-hub-api.service
[Unit]
Description=Health Hub API
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/health-hub
EnvironmentFile=/opt/health-hub/server.env
ExecStart=/usr/bin/node /opt/health-hub/server/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
